# 配置文件和源码/Python监控进程状态并发实现告警/proms.ini
[program:proms]
command=/opt/prometheus/server/prometheus/prometheus
directory=/opt/prometheus/server/prometheus
stdout_logfile=/home/data/logs/prometheus/sever.log
autostart=true
autorestart=true
redirect_stderr=true
user=root
startsecs=3